# workers/restart-scripts/crontab
*/1 * * * * /home/jrrjdev/workers/restart-scripts/workers.sh